# fogbow-scripts/cron.d/fogbow-monitoring-cron-script
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
1 1 * * * ubuntu bash /home/ubuntu/git/fogbow-monitor-scripts/fogbow-scripts/fogbow-monitoring-create-instances.sh > /var/www/downloads.fogbowcloud.org/fogbow-monitoring/fogbow/create_instances_`date +"\%y_\%h_\%d_\%H_\%M"`.log 2>&1
1 1 * * * ubuntu bash /home/ubuntu/git/fogbow-monitor-scripts/fogbow-scripts/fogbow-monitoring-create-network.sh > /var/www/downloads.fogbowcloud.org/fogbow-monitoring/fogbow/create_network_`date +"\%y_\%h_\%d_\%H_\%M"`.log 2>&1
1 1 * * * ubuntu bash /home/ubuntu/git/fogbow-monitor-scripts/fogbow-scripts/fogbow-monitoring-create-storage.sh > /var/www/downloads.fogbowcloud.org/fogbow-monitoring/fogbow/create_storage_`date +"\%y_\%h_\%d_\%H_\%M"`.log 2>&1
